# tardis-webapp/rpm/postinstall
chkconfig tardis on || true
service tardis start || true
